# bash.conf
[program:bash]
command=/bin/bash
redirect_stderr=true
